# edict.service
# 三省六部 · systemd 服务模板
#
# 当前主部署方式为 edict/docker-compose.yml。本服务只负责托管当前 Compose 主栈。
# 部署说明：
#   sudo cp edict.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now edict

[Unit]
Description=Edict 三省六部 Docker Compose Stack
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/edict/edict
ExecStart=/usr/bin/docker compose up -d --build
ExecStop=/usr/bin/docker compose down
ExecReload=/usr/bin/docker compose up -d --build
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
